# packaging/debian-matrix/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop led-matrix.service 2>/dev/null || true
        systemctl disable led-matrix.service 2>/dev/null || true
    fi
fi

exit 0
